# assets/services/pulp-ui/run
#!/usr/bin/execlineb -P
#npm --prefix /pulp-ui/pulp_ui/app/ run start:dev